# conf/supervisor.d/apache2.conf

[program:apache2]
command=/bin/bash -c "/usr/local/bin/apache2-foreground.sh"
user = root
autostart = true